# crontab
* * * * * docker run --rm example_app_image:latest / task.py
